# gunicorn.unix.conf
[program:gunicorn]
command=gunicorn --workers 4 -k gevent hello:app -b unix:/tmp/deploy_app.sock
directory=/deploy/app
