# app/gunicorn.unix.conf
[program:gunicorn]
command=gunicorn --workers 4 -k gevent 'app:create_app()' -b unix:/tmp/deploy_app.sock
directory=/deploy
